# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/etc/init.d/sshd start

[program:postgresql]
command=/bin/sh /usr/scripts/start_postgres.sh
